# etc/profile.d/s11_aliases_general.sh
alias h='     history'
alias H='     history'
alias hgrep=' history | grep '
alias hmore=' history | more '
alias ll='ls -lAtrF'
shopt -s histappend 
